# supervisord.conf
[supervisord]
nodaemon = true

[program:storj-daemon]
command = storjshare daemon -F

[program:storj-share]
command = storj-start
